# debian/dcc-io-daemon.service
[Unit]
Description=DCC IO Daemon - RESTful API for DCC command stations
After=network.target

[Service]
Type=simple
User=dcc-io
Group=dcc-io
WorkingDirectory=/var/lib/dcc-io-daemon
Environment="JMRI_HOME=/opt/jmri/current"
Environment="JAVA_OPTS=-Xmx512m"

# Build classpath with JMRI
# Note: Wildcards in classpath need to be expanded by shell
# JAR_FILENAME_PLACEHOLDER will be replaced during package build
# Fallback: if placeholder not replaced, use wildcard to find JAR
ExecStart=/bin/bash -c 'JAR_FILE="JAR_FILENAME_PLACEHOLDER"; if [ "$$JAR_FILE" = "JAR_FILENAME_PLACEHOLDER" ]; then JAR_FILE=$$(ls /usr/lib/dcc-io-daemon/dcc-io-daemon-*-jar-with-dependencies.jar 2>/dev/null | head -1); fi; CLASSPATH="$$JAR_FILE:/opt/jmri/current/jmri.jar"; if [ -d /opt/jmri/current/lib ]; then for jar in /opt/jmri/current/lib/*.jar; do [ -f "$$jar" ] && CLASSPATH="$$CLASSPATH:$$jar"; done; fi; exec /usr/bin/java -cp "$$CLASSPATH" org.dccio.daemon.DccIoDaemon 9000'

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dcc-io-daemon

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/dcc-io-daemon

# Allow access to USB, serial devices, and hardware interfaces (GPIO, I2C, SPI)
SupplementaryGroups=dialout plugdev gpio i2c spi input

# Allow access to device files (USB, serial ports)
PrivateDevices=false

[Install]
WantedBy=multi-user.target

