# debian/dcc-io-daemon/DEBIAN/prerm
#!/bin/bash
set -e

# Pre-removal script for dcc-io-daemon

# Stop and disable the service before removal
if systemctl is-active dcc-io-daemon >/dev/null 2>&1; then
    echo "Stopping dcc-io-daemon service..."
    systemctl stop dcc-io-daemon
fi

if systemctl is-enabled dcc-io-daemon >/dev/null 2>&1; then
    echo "Disabling dcc-io-daemon service..."
    systemctl disable dcc-io-daemon
fi

# DEBHELPER
exit 0

